# deployment/files/supervisor/crypto-repo.conf
[program:repo_base]
command=bash /srv/CIP/gunicorn_start
user=uconsumer
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/srv/CIP/logs/gunicorn.log